# agent-patrol/examples/schedules/financeqa-production-hourly.service
[Unit]
Description=Agent Patrol FinanceQA hourly production dry-run local report

[Service]
Type=oneshot
WorkingDirectory=/opt/finance_qa/agent-patrol
Environment=PATH=/root/.nvm/versions/node/v22.22.2/bin:/usr/local/bin:/usr/bin:/bin
Environment=AGENT_PATROL_ENV_FILE=/opt/finance_qa/agent-patrol/examples/schedules/financeqa-production-hourly.env
ExecStart=/opt/finance_qa/agent-patrol/examples/schedules/run-financeqa-dry-run.sh
