# agent-patrol/examples/schedules/financeqa-daily.cron.example
SHELL=/bin/bash
AGENT_PATROL_ROOT=/opt/finance_qa/agent-patrol

# Run two low-frequency local-report-only FinanceQA dry-runs per day.
# Copy financeqa-daily.env.example to financeqa-daily.env and adjust paths before installing.
20 9,17 * * * cd "$AGENT_PATROL_ROOT" && AGENT_PATROL_ENV_FILE=examples/schedules/financeqa-daily.env examples/schedules/run-financeqa-dry-run.sh
