# agent-patrol/examples/schedules/financeqa-daily.timer
[Unit]
Description=Run Agent Patrol FinanceQA daily local report

[Timer]
OnCalendar=*-*-* 09,13,18:10:00
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target

