# agent-patrol/examples/schedules/financeqa-production-hourly.timer
[Unit]
Description=Run Agent Patrol FinanceQA hourly production dry-run local report

[Timer]
OnCalendar=*-*-* *:07:00
Persistent=false
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
