# deploy/systemd/financeqa-ocr-worker.service
[Unit]
Description=FinanceQA Gemini OCR pending worker
Documentation=file:%h/finance_qa/README.md
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=%h/finance_qa
Environment=FINANCEQA_ENV_FILE=%h/finance_qa/.env
ExecStart=%h/finance_qa/bin/financeqa ocr process-pending
TimeoutStartSec=30min
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
SyslogIdentifier=financeqa-ocr-worker
