# deploy/systemd/financeqa-ocr-worker.service
[Unit]
Description=FinanceQA Gemini OCR pending worker
Documentation=file:/root/finance_qa/README.md
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/root/finance_qa
EnvironmentFile=/root/finance_qa/.env
ExecStart=/root/finance_qa/financeqa ocr process-pending
TimeoutStartSec=30min
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
SyslogIdentifier=financeqa-ocr-worker
