# deploy/systemd/financeqa-ocr-worker.timer
[Unit]
Description=Run FinanceQA Gemini OCR worker every 5 minutes
Documentation=/root/finance_qa/README.md

[Timer]
OnUnitInactiveSec=5min
AccuracySec=30s
Unit=financeqa-ocr-worker.service

[Install]
WantedBy=timers.target
